# docker/files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/etc/init.d/php5-fpm start

[program:nginx]
command=/etc/init.d/nginx start
